# ansible/roles/hostconfig/files/systemd/otvl_default_route.service
[Unit]
Description=otvl_default_route script
After=network.target local-fs.target

[Service]
ExecStart=sh -c /usr/local/bin/otvl_default_route.sh
Restart=always
RestartSec=60

[Install]
WantedBy=multi-user.target
